# ansible/files/gandalf_config/gandalf_salt.service
[Unit]
Description=Configuration management via gandalf salt
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/gandalf/gandalf --config /opt/gandalf/config.yaml salt --apply --store-result
User=gandalf
Group=gandalf
StandardOutput=journal
StandardError=journal
StateDirectory=gandalf
CacheDirectory=gandalf

# Security hardening
PrivateTmp=true

[Install]
WantedBy=multi-user.target
